# fs.ini
[program:fs]
command=/usr/lib/jvm/java-8-openjdk-amd64 -jar client.jar -b
umask=022
user=root
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
stopsignal=QUIT
killasgroup=true